# zeek_docker/Zeek-source-cluster/supervisor-slave.conf
[supervisord]
nodaemon=true

[program:bro]
command=/opt/bro/bin/broctl deploy
startretries=3

[program:sshd]
command=/usr/sbin/sshd -D
